# rh-containers-x-container-2/config/chrome-direct.conf
; Priority 0 - fix-permissions (root), 1 - xvfb & fluxbox, 5 - x11vnc, 10 - noVNC, 15 - chrome direct

[program:fix-permissions]
priority=0
command=/opt/bin/fix-permissions.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
user=root
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xvfb]
priority=0
command=/opt/bin/start-xvfb.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
redirect_stderr=false
stdout_logfile=/var/log/supervisor/xvfb-stdout.log
stderr_logfile=/var/log/supervisor/xvfb-stderr.log

[program:vnc]
priority=5
command=/opt/bin/start-vnc.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
redirect_stderr=false
stdout_logfile=/var/log/supervisor/vnc-stdout.log
stderr_logfile=/var/log/supervisor/vnc-stderr.log

[program:novnc]
priority=10
command=/opt/bin/start-novnc.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
redirect_stderr=false
stdout_logfile=/var/log/supervisor/novnc-stdout.log
stderr_logfile=/var/log/supervisor/novnc-stderr.log

[program:chrome-direct]
priority=15
command=/opt/bin/start-chrome-direct.sh
stopasgroup=true
autostart=true
autorestart=true
startsecs=3
startretries=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:ws-server]
priority=20
command=python3 /opt/x-automation/api/ws_server.py
directory=/opt/x-automation
autostart=true
autorestart=true
startsecs=2
startretries=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=CDP_URL="http://127.0.0.1:9222",WS_LISTEN_PORT="18000",X_AUTOMATION_ALLOW_MUTATIONS="%(ENV_X_AUTOMATION_ALLOW_MUTATIONS)s",CHROME_MCP_BIN="/usr/lib/node_modules/chrome-devtools-mcp/build/src/bin/chrome-devtools-mcp.js",NODE_BIN="/usr/bin/node"
